# sensible-cvsb983.conf
# sensible-cvsb983
#

description	"sensible-cvsb983 daemon"
author		"Andrew Thornton"

start on (starting xbmc)
stop on (xbmc-do-stop or runlevel [!2345])

respawn

exec /usr/local/bin/sensible-cvsb983
